# deploy/agentbot.service
[Unit]
Description=Agent Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=agent
WorkingDirectory=/opt/agentbot
ExecStart=/opt/agentbot/.venv/bin/python -m app.main
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
